# llm-proxy.service
[Unit]
Description=LLM Proxy FastAPI Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/github-runner/action-runner/llm-proxy
ExecStart=/home/github-runner/action-runner/llm-proxy/.venv/bin/python -m uvicorn proxy:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3
User=github-runner
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target